# systemd/epipylon-apt-update.service
[Unit]
Description=Update system software packages

[Service]
Type=oneshot
ExecStart=/usr/sbin/epipylon-apt-update
StandardError=syslog
